# smc-signal-tracker.service
[Unit]
Description=TradingExpert Signal Tracker Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/evans/smc-scalp-signals
Environment="PYTHONPATH=/home/evans/smc-scalp-signals"
ExecStart=/usr/bin/python3 /home/evans/smc-scalp-signals/signal_tracker.py
Restart=always
RestartSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=smc-tracker

[Install]
WantedBy=multi-user.target
