# smc-health-watchdog.service
[Unit]
Description=SMC Institutional Health Watchdog
After=network.target smc-signal-service.service smc-admin-dashboard.service

[Service]
Type=simple
User=evans
WorkingDirectory=/home/evans/smc-scalp-signals
ExecStart=/home/evans/smc-scalp-signals/venv/bin/python monitoring/health_check.py
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
