# smc-interactive-bot.service
[Unit]
Description=Interactive Telegram Bot for TradingExpert
After=network.target

[Service]
Type=simple
User=evans
WorkingDirectory=/home/evans/smc-scalp-signals
Environment="PYTHONPATH=/home/evans/smc-scalp-signals"
ExecStart=/home/evans/smc-scalp-signals/venv/bin/python /home/evans/smc-scalp-signals/app/interactive_bot.py
Restart=always
RestartSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=smc-bot

[Install]
WantedBy=multi-user.target
