# smc-signal-service.service
[Unit]
Description=Pure Quant Trading Signal Service
After=network.target

[Service]
Type=simple
User=evans
WorkingDirectory=/home/evans/smc-scalp-signals
Environment="PYTHONPATH=/home/evans/smc-scalp-signals"
ExecStart=/home/evans/smc-scalp-signals/venv/bin/python /home/evans/smc-scalp-signals/signal_service.py
Restart=always
RestartSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=quant-signals

[Install]
WantedBy=multi-user.target
