# smc-admin-dashboard.service
[Unit]
Description=TradingExpert Admin Dashboard Server
After=network.target

[Service]
Type=simple
User=evans
WorkingDirectory=/home/evans/smc-scalp-signals
Environment="PYTHONPATH=/home/evans/smc-scalp-signals"
ExecStart=/home/evans/smc-scalp-signals/venv/bin/python /home/evans/smc-scalp-signals/admin_server.py
Restart=always
RestartSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=smc-admin

[Install]
WantedBy=multi-user.target
