# monitoring/smc-watchdog.service
[Unit]
Description=Signal Service Watchdog (Health Checks)
After=network.target

[Service]
Type=oneshot
User=evans
WorkingDirectory=/home/evans/smc-scalp-signals
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/evans/smc-scalp-signals/venv/bin/python /home/evans/smc-scalp-signals/monitoring/watchdog.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=smc-watchdog
