# monitoring/smc-daily-report.service
[Unit]
Description=Daily Signal Service Performance Report
After=network.target

[Service]
Type=oneshot
User=evans
WorkingDirectory=/home/evans/smc-scalp-signals
Environment="PYTHONUNBUFFERED=1"
ExecStart=/home/evans/smc-scalp-signals/venv/bin/python /home/evans/smc-scalp-signals/monitoring/daily_report.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=smc-daily-report
